# home/.config/systemd/user/plasma-kwin_wayland.service.d/override.conf
[Service]
Environment=KWIN_DRM_DEVICES=/dev/dri/card1